# 02-siem-configs/wazuh/custom-rules.xml
<!-- 
  Wazuh Custom Detection Rules
  Mini-SOC Project
  
  Installation:
  1. Copy to: /var/ossec/etc/rules/local_rules.xml
  2. Restart Wazuh: sudo systemctl restart wazuh-manager
  3. Verify: /var/ossec/bin/ossec-logtest
-->

<group name="custom,">

  <!-- ============================================ -->
  <!-- USE CASE 1: Brute-Force Attack Detection    -->
  <!-- ============================================ -->
  
  <!-- Rule 100001: Multiple failed logins (Windows) -->
  <rule id="100001" level="10" frequency="5" timeframe="120">
    <if_matched_sid>60122</if_matched_sid>
    <description>Multiple Windows failed login attempts (Possible brute-force)</description>
    <mitre>
      <id>T1110.001</id>
    </mitre>
    <group>authentication_failures,windows,brute_force,</group>
  </rule>

  <!-- Rule 100002: Multiple failed logins (Linux SSH) -->
  <rule id="100002" level="10" frequency="5" timeframe="120">
    <if_matched_sid>5551</if_matched_sid>
    <description>Multiple SSH failed login attempts (Possible brute-force)</description>
    <mitre>
      <id>T1110.001</id>
    </mitre>
    <group>authentication_failures,ssh,brute_force,</group>
  </rule>

  <!-- Rule 100003: Successful login after brute-force -->
  <rule id="100003" level="12">
    <if_matched_sid>100001,100002</if_matched_sid>
    <if_sid>60103,5501</if_sid>
    <description>Successful login after brute-force attempt (Potential breach)</description>
    <mitre>
      <id>T1110.001</id>
    </mitre>
    <group>authentication_success,brute_force,critical,</group>
  </rule>

  <!-- ============================================ -->
  <!-- USE CASE 2: Admin Privilege Abuse           -->
  <!-- ============================================ -->

  <!-- Rule 100010: Privileged account login outside business hours -->
  <rule id="100010" level="8">
    <if_sid>60103</if_sid>
    <match>Administrator|admin|root</match>
    <time>6 pm - 6 am</time>
    <description>Administrator login outside business hours</description>
    <mitre>
      <id>T1078.003</id>
    </mitre>
    <group>authentication,privileged_account,</group>
  </rule>

  <!-- Rule 100011: Sensitive file access by admin -->
  <rule id="100011" level="7">
    <if_sid>554</if_sid>
    <match>/etc/shadow|/etc/passwd|SAM|SYSTEM</match>
    <description>Administrator accessed sensitive authentication files</description>
    <mitre>
      <id>T1003.008</id>
    </mitre>
    <group>file_access,credential_access,</group>
  </rule>

  <!-- Rule 100012: Multiple privilege escalations -->
  <rule id="100012" level="10" frequency="3" timeframe="300">
    <if_matched_sid>4673</if_matched_sid>
    <description>Multiple privilege escalation attempts detected</description>
    <mitre>
      <id>T1078</id>
    </mitre>
    <group>privilege_escalation,suspicious_activity,</group>
  </rule>

  <!-- ============================================ -->
  <!-- USE CASE 3: Web Application Attacks         -->
  <!-- ============================================ -->

  <!-- Rule 100020: SQL Injection attempt -->
  <rule id="100020" level="10">
    <if_sid>31100,31101,31102,31103</if_sid>
    <url>select|union|insert|update|delete|drop|exec|script</url>
    <description>Possible SQL injection attack detected</description>
    <mitre>
      <id>T1190</id>
    </mitre>
    <group>web_attack,sql_injection,</group>
  </rule>

  <!-- Rule 100021: Cross-Site Scripting (XSS) attempt -->
  <rule id="100021" level="9">
    <if_sid>31100,31101,31102,31103</if_sid>
    <url>&lt;script|javascript:|onerror=|onload=</url>
    <description>Possible XSS attack detected</description>
    <mitre>
      <id>T1059.007</id>
    </mitre>
    <group>web_attack,xss,</group>
  </rule>

  <!-- Rule 100022: Directory traversal attempt -->
  <rule id="100022" level="9">
    <if_sid>31100,31101,31102,31103</if_sid>
    <url>../|..%2F|..%5c</url>
    <description>Directory traversal attack detected</description>
    <mitre>
      <id>T1190</id>
    </mitre>
    <group>web_attack,directory_traversal,</group>
  </rule>

  <!-- Rule 100023: Multiple 404 errors (Scanning) -->
  <rule id="100023" level="8" frequency="10" timeframe="60">
    <if_matched_sid>31303</if_matched_sid>
    <description>Multiple 404 errors - Possible web scanning activity</description>
    <mitre>
      <id>T1595.002</id>
    </mitre>
    <group>web_scan,reconnaissance,</group>
  </rule>

  <!-- ============================================ -->
  <!-- USE CASE 4: DNS Exfiltration                -->
  <!-- ============================================ -->

  <!-- Rule 100030: Unusually long DNS query -->
  <rule id="100030" level="9">
    <if_sid>86600</if_sid>
    <regex type="pcre2">query[^\s]{60,}</regex>
    <description>Unusually long DNS query detected (Possible exfiltration)</description>
    <mitre>
      <id>T1071.004</id>
    </mitre>
    <group>dns,exfiltration,</group>
  </rule>

  <!-- Rule 100031: High frequency DNS queries -->
  <rule id="100031" level="8" frequency="50" timeframe="60">
    <if_matched_sid>86600</if_matched_sid>
    <same_source_ip />
    <description>High volume DNS queries from single source (Possible C2 or exfiltration)</description>
    <mitre>
      <id>T1071.004</id>
    </mitre>
    <group>dns,exfiltration,command_control,</group>
  </rule>

  <!-- Rule 100032: DNS query to suspicious TLD -->
  <rule id="100032" level="7">
    <if_sid>86600</if_sid>
    <match>.tk|.ml|.ga|.cf|.gq</match>
    <description>DNS query to suspicious free domain TLD</description>
    <mitre>
      <id>T1071.004</id>
    </mitre>
    <group>dns,suspicious_domain,</group>
  </rule>

  <!-- ============================================ -->
  <!-- USE CASE 5: Persistence Mechanisms          -->
  <!-- ============================================ -->

  <!-- Rule 100040: Registry Run key modification (Windows) -->
  <rule id="100040" level="10">
    <if_sid>61600</if_sid>
    <regex type="pcre2">HKLM\\Software\\Microsoft\\Windows\\CurrentVersion\\Run|HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run</regex>
    <description>Registry Run key modified - Possible persistence mechanism</description>
    <mitre>
      <id>T1547.001</id>
    </mitre>
    <group>persistence,registry,windows,</group>
  </rule>

  <!-- Rule 100041: Startup folder file creation -->
  <rule id="100041" level="9">
    <if_sid>61603</if_sid>
    <regex type="pcre2">Start Menu\\Programs\\Startup</regex>
    <description>File created in Startup folder - Possible persistence</description>
    <mitre>
      <id>T1547.001</id>
    </mitre>
    <group>persistence,startup,windows,</group>
  </rule>

  <!-- Rule 100042: Scheduled task creation -->
  <rule id="100042" level="9">
    <if_sid>106</if_sid>
    <match>schtasks|at.exe|Task Scheduler</match>
    <description>Scheduled task created - Possible persistence mechanism</description>
    <mitre>
      <id>T1053.005</id>
    </mitre>
    <group>persistence,scheduled_task,windows,</group>
  </rule>

  <!-- Rule 100043: Cron job modification (Linux) -->
  <rule id="100043" level="9">
    <if_sid>2904</if_sid>
    <match>/etc/cron|/var/spool/cron</match>
    <description>Cron job modified - Possible persistence mechanism</description>
    <mitre>
      <id>T1053.003</id>
    </mitre>
    <group>persistence,cron,linux,</group>
  </rule>

  <!-- ============================================ -->
  <!-- USE CASE 6: Lateral Movement                -->
  <!-- ============================================ -->

  <!-- Rule 100050: Multiple SMB connections to different hosts -->
  <rule id="100050" level="9" frequency="5" timeframe="300">
    <if_matched_sid>60000</if_matched_sid>
    <match>SMB|445</match>
    <description>Multiple SMB connections to different hosts (Possible lateral movement)</description>
    <mitre>
      <id>T1021.002</id>
    </mitre>
    <group>lateral_movement,smb,network,</group>
  </rule>

  <!-- Rule 100051: PsExec usage detected -->
  <rule id="100051" level="12">
    <if_sid>106</if_sid>
    <match>psexec|paexec</match>
    <description>PsExec execution detected - Possible lateral movement</description>
    <mitre>
      <id>T1570</id>
    </mitre>
    <group>lateral_movement,psexec,critical,</group>
  </rule>

  <!-- Rule 100052: WMI remote execution -->
  <rule id="100052" level="10">
    <if_sid>106</if_sid>
    <match>wmic.exe</match>
    <regex type="pcre2">/node:</regex>
    <description>WMI remote execution detected - Possible lateral movement</description>
    <mitre>
      <id>T1047</id>
    </mitre>
    <group>lateral_movement,wmi,</group>
  </rule>

  <!-- Rule 100053: RDP connections from unusual sources -->
  <rule id="100053" level="8">
    <if_sid>60103</if_sid>
    <match>EventID: 4624</match>
    <regex type="pcre2">LogonType: 10</regex>
    <description>RDP login detected - Monitor for lateral movement</description>
    <mitre>
      <id>T1021.001</id>
    </mitre>
    <group>lateral_movement,rdp,</group>
  </rule>

  <!-- ============================================ -->
  <!-- ADDITIONAL DETECTION RULES                  -->
  <!-- ============================================ -->

  <!-- Rule 100060: Mimikatz usage detected -->
  <rule id="100060" level="15">
    <if_sid>106</if_sid>
    <match>mimikatz|sekurlsa|kerberos::golden</match>
    <description>Mimikatz credential dumping tool detected - CRITICAL</description>
    <mitre>
      <id>T1003.001</id>
    </mitre>
    <group>credential_access,mimikatz,critical,</group>
  </rule>

  <!-- Rule 100061: PowerShell obfuscation detected -->
  <rule id="100061" level="10">
    <if_sid>91816</if_sid>
    <regex type="pcre2">-encodedcommand|-enc|-e\s+[A-Za-z0-9+/=]{50,}</regex>
    <description>Obfuscated PowerShell command detected</description>
    <mitre>
      <id>T1059.001</id>
    </mitre>
    <group>execution,powershell,obfuscation,</group>
  </rule>

  <!-- Rule 100062: Suspicious process spawned by Office application -->
  <rule id="100062" level="12">
    <if_sid>106</if_sid>
    <regex type="pcre2">WINWORD\.EXE|EXCEL\.EXE|POWERPNT\.EXE</regex>
    <match>cmd.exe|powershell.exe|wscript.exe|cscript.exe</match>
    <description>Office application spawned suspicious process - Possible malware</description>
    <mitre>
      <id>T1566.001</id>
    </mitre>
    <group>execution,office,malware,critical,</group>
  </rule>

</group>
